# addon-icecast/rootfs/etc/s6-overlay/s6-rc.d/icecast/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Icecast
# Runs icecast script
# ==============================================================================

bashio::log.info "Starting Icecast..."

exec /usr/bin/icecast.sh
